# supervisord.conf
[supervisord]

[program:inbox-celery]
command = /home/ubuntu/casepro/env/bin/celery worker -A casepro -Q celery,sync -B -n casepro.celery --loglevel=ERROR --autoscale=2,1
autostart=true